# vendor/supervisord_pd-emitter.conf
# deploy to /etc/supervisor/conf.d/ (In case of Ubuntu)
[supervisord]

[program:pd-emitter]
directory = /opt/pd-emitter/
environment=RUN_ENV="production"
command = bundle exec rake start
# if using rbenv, `bundle exec ...` replace to `/home/USER/.rbenv/shims/bundle exec ...`. see README.md#daemonize
redirect_stderr = true

[group:main]
programs = pd-emitter
priority = 999

