# vendor/supervisord_emitter.conf
# deploy to /etc/supervisor/conf.d/ (case of Ubuntu)
[supervisord]
# depend on run env.
#nodaemon = true

[program:emitter]
directory = /opt/emitter/
environment=RUN_ENV="production"
command = bundle exec rake boot
# if using rbenv, `bundle exec ...` replace to `/home/USER/.rbenv/shims/bundle exec ...`. see README.md#daemonize
redirect_stderr = true

[group:main]
programs = emitter
priority = 999

